# docker/supervisord/influxd.conf
[program:influxd]
command=/usr/bin/influxd -config /etc/influxdb/influxdb.conf
autostart=true
autorestart=true
user=influxdb
