# docker/supervisord/mongod.conf
[program:mongod]
command=/usr/bin/mongod --config /etc/mongod.conf
autostart=true
autorestart=true
user=mongodb
